# frontend/supervisor/supervisord.prod.conf
; ── Supervisor — Production ───────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=warn

; ── PHP-FPM ───────────────────────────────────────────────────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── RabbitMQ Consumer ────────────────────────────────────────────────────────
; In prod: longer time-limit, more restarts allowed.
[program:rabbitmq-consumer]
command=/var/www/html/vendor/bin/drush --root=/var/www/html rabbitmq:consume --time-limit=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
priority=20
startretries=10
startsecs=10
; Back off on repeated crashes: wait up to 60s between restarts
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
